# deploy/com.garmin-coach.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.garmin-coach.web</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/anthonytrieu/Desktop/garmin/.venv/bin/garmin-coach-web</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/anthonytrieu/Desktop/garmin</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/anthonytrieu/Library/Logs/garmin-coach.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/anthonytrieu/Library/Logs/garmin-coach.log</string>
</dict>
</plist>
